# gateway_deb/etc/systemd/system/gateway.service
[Unit]
Description = gateway

[Service]
Type = simple
ExecStart = /home/gateway/gateway.sh
Restart = always


[Install]
WantedBy = multi-user.target