# gateway_deb/DEBIAN/postrm
#!/bin/bash

echo "......................................................"
echo "gateway remove"
echo "......................................................"

systemctl disable gateway